# data/com.redhat.spice.vdagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>KeepAlive</key>
		<dict>
			<key>Crashed</key>
			<false/>
		</dict>
		<key>Label</key>
		<string>com.redhat.spice.vdagent</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/local/bin/spice-vdagent</string>
			<string>-x</string>
			<string>-S</string>
			<string>/var/run/spice-vdagent-sock</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
	</dict>
</plist>
