# data/spice-vdagentd.socket
[Unit]
Description=Activation socket for spice guest agent daemon

[Socket]
ListenStream=/var/run/spice-vdagentd/spice-vdagent-sock
